# etc/profile.d/vars.sh
export VISUAL=kak
export EDITOR=kak
export QT_QPA_PLATFORM=wayland
export SDL_VIDEODRIVER=wayland
export MOZ_ENABLE_WAYLAND=1
export DOTNET_ROOT=$HOME/dotnet
export PATH=$PATH:$HOME/dotnet
export CLR_OPENSSL_VERSION_OVERRIDE=47
