# etc/profile.d/vars.sh
export VISUAL=nvim
export EDITOR=nvim
export QT_QPA_PLATFORM=wayland
export SDL_VIDEODRIVER=wayland
export MOZ_ENABLE_WAYLAND=1
export DOTNET_ROOT=$HOME/dotnet
export PATH=$PATH:$HOME/dotnet
export PATH=$PATH:$HOME/.local/bin
export CLR_OPENSSL_VERSION_OVERRIDE=47
